# edge/systemd/petir-sync.service.tmpl
[Unit]
Description=Petir edge sync (single-run push to lab server)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__EDGE_DIR__
ExecStart=__VENV_PYTHON__ -m sync_worker.run
EnvironmentFile=__ENV_FILE__
TimeoutStartSec=60
